# pkg/post-install.sh
#!/bin/bash
systemctl enable aerospike-prometheus-exporter.service